# Script-inits/test_server-init.service
[Unit]
Description=start nodejs server

[Service]
ExecStart=/root/node_modules/ws/test_server-init.sh

[Install]
WantedBy=multi-user.target